# crontab
@reboot ./hello.sh
*/10 * * * * ./goodbye.sh
